# etc/supervisord.d/jmxtrans.ini
[program:jmxtrans]
command=/opt/qnib/jmxtrans/bin/start.sh
redirect_stderr=true
#stdout_logfile=syslog
stdout_logfile=/var/log/supervisor/%(program_name)s.log
startsecs=1
exitcodes=0
