# etc/supervisord.d/kafka.ini
[program:kafka]
stdout_logfile=/var/log/supervisor/%(program_name)s.log
command=/opt/hello-samza/deploy/kafka/bin/kafka-server-start.sh /opt/hello-samza/deploy/kafka/config/server.properties
redirect_stderr=true
exitcodes=0
startsecs=10
priority=500
