# etc/supervisord.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
redirect_stdery=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
autostart=true
startsecs=5
exitcodes=0
