# etc/supervisord.d/yarn-nodemanager.ini
[program:yarn-nodemanager]
stdout_logfile=/var/log/supervisor/%(program_name)s.log
command=/opt/hello-samza/deploy/yarn/bin/yarn nodemanager
redirect_stderr=true
exitcodes=0
startsecs=10
priority=999
