# etc/supervisord.d/zookeeper.ini
[program:zookeeper]
command=/opt/hello-samza/deploy/zookeeper/bin/zkServer.sh start-foreground
redirect_stderr=true
exitcodes=0
startsecs=10
priority=200
stdout_logfile=/var/log/supervisor/%(program_name)s.log
